# cmake/hello_service/systemd/hello.service
[Unit]
Description=hello service

[Service]
ExecStart=/usr/bdk/sbin/hello
Restart=always
Type=simple

[Install]
WantedBy=multi-user.target
